# _system/systemd/tradewave-rt-backup.timer
[Unit]
Description=Daily TradeWave RT DB backup

[Timer]
OnCalendar=*-*-* 08:30:00
Persistent=true

[Install]
WantedBy=timers.target
