# tradewave_realtime/deploy/units/tradewave-rt-levels.service
[Unit]
Description=TradeWave RT - levels sync (keyprovider captures + hourly EODHD + map rebuild)
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/flask/tradewave_realtime
EnvironmentFile=/etc/tradewave_realtime/secrets.env
ExecStart=/home/flask/tradewave_realtime/.venv/bin/flask levels-sync
